# Systemd-Services/hyprsunset-off.service
[Unit]
Description=Turn off hyprsunset

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill hyprsunset
