# Systemd-Services/hyprsunset-on.service
[Unit]
Description=Turn on hyprsunset

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprsunset --temperature 3500
