# Systemd-Services/wlsunset-off.service
[Unit]
Description=Turn off wlsunset

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill wlsunset
