# systemd-services/wlsunset-on.service
[Unit]
Description=Turn on wlsunset

[Service]
Type=oneshot
ExecStart=/usr/bin/wlsunset -T 3501 -t 3500
